# systemd/runner-install.service
[Unit]
Description=runner first-boot installer
Wants=network-online.target
After=network-online.target

StartLimitIntervalSec=0

# Don't rerun once installed.
ConditionPathExists=!/var/lib/runner/installed

[Service]
Type=oneshot
EnvironmentFile=-/etc/runner/config.env
ExecStart=/usr/local/lib/runner/bootstrap.sh
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
