# systemd/template-appliance-secondary.service
[Unit]
Description=template-appliance secondary mode
Wants=network.target
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/template-appliance/config.env
ExecStart=/usr/local/lib/template-appliance/secondary-mode.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
