# systemd/template-appliance-healthcheck.service
[Unit]
Description=template-appliance healthcheck/failover

[Service]
Type=oneshot
EnvironmentFile=-/etc/template-appliance/config.env
ExecStart=/usr/local/lib/template-appliance/healthcheck.sh
